# ansible/roles/aegir/templates/etc/systemd/system/mysql.service.d/filelimit.conf
# {{ ansible_managed }}
[Service]
LimitNOFILE=1024000
